# ticketbot.service
# /etc/systemd/system/ticketbot.service
#
# systemd-Unit für den MSK Scripts Discord Ticket-Bot
# ────────────────────────────────────────────────────
# Installation:
#   sudo cp ticketbot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ticketbot.service
#
# Nützliche Befehle:
#   sudo systemctl status  ticketbot.service   # Status anzeigen
#   sudo systemctl restart ticketbot.service   # Neustart
#   sudo systemctl stop    ticketbot.service   # Stoppen
#   sudo journalctl -u ticketbot.service -f    # Live-Logs (journald)
#   sudo journalctl -u ticketbot.service -f --output=cat   # Follow live logs with colors

[Unit]
Description=MSK Scripts Discord Ticket-Bot
Documentation=https://github.com/MSK-Scripts
# Erst starten, wenn das Netzwerk vollständig verfügbar ist
After=network-online.target
Wants=network-online.target
# Neustart-Limit: max 5 Versuche in 120 Sekunden
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
# ── Benutzer & Pfade ─────────────────────────────────────────────────────────
Type=simple
User=discord
Group=discord
WorkingDirectory=/opt/discord_ticketbot

# ── Umgebung ─────────────────────────────────────────────────────────────────
# .env-Datei wird von dotenv automatisch geladen;
# für systemd zusätzlich hier referenzieren (doppelte Absicherung).
EnvironmentFile=/opt/discord_ticketbot/.env

# Node.js direkt als Interpreter, index.js als Einstiegspunkt
ExecStart=/usr/bin/node /opt/discord_ticketbot/index.js

# ── Restart-Verhalten ─────────────────────────────────────────────────────────
# Neustart bei jedem unerwarteten Exit (Absturz, OOM, etc.)
Restart=on-failure
RestartSec=10s

# ── Logging ───────────────────────────────────────────────────────────────────
# stdout/stderr an journald weiterleiten (journalctl -u ticketbot abrufbar)
# ANSI-Farben aktivieren: --output=cat beim Lesen oder SYSTEMD_COLORS=1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ticketbot
Environment=FORCE_COLOR=1

# ── Ressourcen & Sicherheit ───────────────────────────────────────────────────
# Verhindert, dass der Bot Root-Rechte erlangen kann
NoNewPrivileges=true
# Schränkt sichtbares Filesystem ein
PrivateTmp=true
ProtectSystem=strict
# Erlaubt Schreibzugriff nur auf das Projektverzeichnis (Logs, data/)
ReadWritePaths=/opt/discord_ticketbot

# Prozess nicht beenden, wenn er zu viel Speicher braucht (OOM-Killer)
OOMScoreAdjust=-500

# ── Graceful Shutdown ─────────────────────────────────────────────────────────
# Gibt dem Bot 30 Sekunden für sauberes Beenden, bevor SIGKILL gesendet wird
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
